# hermes-office.service
[Unit]
Description=Hermes 3D Office Dashboard
After=network.target

[Service]
Type=simple
User=j1admin
WorkingDirectory=/home/j1admin/hermes-3d-office
ExecStart=/usr/bin/python3 /home/j1admin/hermes-3d-office/server.py
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
